# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx-ix.conf
user=root
numprocs=1
autostart=true
autorestart=true
