# etc/supervisor/conf.d/postgresql.conf
[program:postgres]
command=/etc/init.d/postgresql start
user=root
stdout_logfile=/var/log/pgsql.log
stderr_logfile=/var/log/pgsql.log
autostart=true
autorestart=false
